# debian/postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
	update-rc.d scribe remove >/dev/null
fi

exit 0
